# etc/rc.local
#!/bin/bash

# Turn USB power off
sudo uhubctl --action off -l 2 --repeat 100 > /dev/null

