# etc/update-motd.d/55-motd
#!/bin/sh

cat /etc/update-motd.d/motd.txt
